# os/systemd/k3s/k3s-kubeconfig-setup.service
[Unit]
Description=Setup K3s kubeconfig for users
Documentation=https://k3s.io
After=k3s.service
Wants=k3s.service

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
ExecStart=/usr/local/bin/setup-k3s-kubeconfig.sh
Restart=on-failure
RestartSec=30 