# os/systemd/k3s/k3s-load-images.service
[Unit]
Description=Load K3s embedded container images for offline operation
Documentation=https://k3s.io
Before=k3s.service
After=network-online.target
Wants=network-online.target

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
ExecStart=/usr/local/bin/k3s-load-images.sh
TimeoutStartSec=300
Restart=no
StandardOutput=journal
StandardError=journal 