# os/systemd/microshift-kubeconfig-setup.service
[Unit]
Description=MicroShift Kubeconfig Setup
After=microshift.service
Wants=microshift.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup-microshift-kubeconfig.sh
RemainAfterExit=yes
User=root

[Install]
WantedBy=multi-user.target 