# os/systemd/microshift-kubeconfig-setup.service
[Unit]
Description=Setup kubeconfig for MicroShift
After=microshift.service
Wants=microshift.service
Requisite=microshift.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup-kubeconfig.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
User=root

[Install]
WantedBy=multi-user.target 