# os/systemd/setup-offline-images.service
[Unit]
Description=Setup Offline Container Images for MicroShift
Documentation=https://github.com/ramaedge/os-builder
Before=microshift.service crio.service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup-offline-images.sh
RemainAfterExit=true
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target 